# iems5722_app.conf
[program:assignment_3]
command = /home/ubuntu/IEMS5722-A3/venv/bin/gunicorn -b 0.0.0.0:8000 -w 4 assignment_3:app
directory = /home/ubuntu/IEMS5722-A3
user = ubuntu
autostart = true
autorestart = true
stdout_logfile = /home/ubuntu/IEMS5722-A3/app.log
redirect_stderr = true
